# deploy/systemd/telegram-bot-web.service
[Unit]
Description=Telegram OAuth web service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/telegram-bot
EnvironmentFile=-/etc/telegram-bot.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m app.infra.google_oauth_server
Restart=on-failure
RestartSec=5
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
